# deploy/systemd/trakshya-dashboard.service
[Unit]
Description=TRAKSHYA WAF Dashboard
Documentation=https://github.com/Pradyu12/KALKI-WAF
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/trakshya-waf
Environment=NODE_ENV=production
EnvironmentFile=/opt/trakshya-waf/.env
ExecStart=node server.js
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
