# deploy/systemd/trakshya-api.service
[Unit]
Description=TRAKSHYA WAF Management API
Documentation=https://github.com/Pradyu12/KALKI-WAF
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/trakshya-waf
EnvironmentFile=/opt/trakshya-waf/.env
ExecStart=/opt/trakshya-waf/build/trakshya-api
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
